# services/power-control-wake.service
[Unit]
Description=Power Control Wake Service
Documentation=https://github.com/username/linux-power-manager
After=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/power-control.sh wake
ExecStartPost=/usr/local/bin/wifi-intel-optimizer.sh runtime
RemainAfterExit=no
TimeoutStartSec=20
User=root
Environment=DISPLAY=:0

[Install]
WantedBy=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target
